# etc/inittab
::sysinit:/etc/init.d/rcS
tty1::respawn:/sbin/getty 9600 tty1
tty2::respawn:/sbin/getty 9600 tty2
#tty3::respawn:/sbin/getty 9600 tty3
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/sync
::shutdown:/bin/umount -a -r

# remove '#' of next line to support serial console logins
# ttyS0::respawn:/sbin/getty -L 9600 ttyS0

